# Thesis/code/browser_boot.conf
start on desktop-start
stop on desktop-end

exec chromium-browser --use-gl=egl --kiosk http://solarpoweredbikes.tudelft.nl/eud/app_status