# Thesis/code/ODROID_booter.conf
start on startup
task
exec /home/odroid/Documents/SUNRISE-Odroid/Bash/ODROID_booter.sh