# Migration-files/upstart/browser_boot.conf
start on desktop-start
stop on desktop-end

exec chromium-browser --use-gl=egl