# Migration-files/upstart/xvnc_killer.conf
start on startup
task
exec /home/odroid/Documents/SUNRISE-Odroid/Bash/xvnc_killer.sh