# Bash/browser_boot.conf
start on desktop-start
stop on desktop-end

exec chromium-browser --use-gl=egl --noerrdialogs --disable-session-crashed-bubble --disable-infobars --kiosk http://localhost:8000/Local_webpage