# Bash/firefox.conf
start on desktop-start
stop on desktop-end

exec firefox