# Bash/Local_server_booter.conf
start on desktop-start
stop on desktop-end

exec /home/odroid/Documents/SUNRISE-Odroid/Bash/Local_server_booter.sh